# bin/recollecing.service
[Unit]
Description=The recollecing service.

[Service]
ExecStart=/home/pi/.local/bin/reco --db-file /media/usb/db.sqlite run
Restart=on-failure

[Install]
WantedBy=default.target
